# src/service/msi-ec-service.service.in
[Unit]
Description=Msi EC Control Service
After=network.target

[Service]
Type=simple
ExecStartPre=/sbin/modprobe msiecmodule
ExecStart=@MSICONTROLLER_SERVICE_EXEC@
Restart=on-failure
User=root

[Install]
WantedBy=multi-user.target
